# solab-bot.service
[Unit]
Description=SoLab Telegram Bot
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/projects/solab
Environment=PYTHONPATH=/root/projects/solab
ExecStart=/root/projects/solab/venv/bin/python /root/projects/solab/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
