# solab-bot.service
[Unit]
Description=SoLab Telegram Bot
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/root/telegram-bot/solab
Environment=PYTHONPATH=/home/root/telegram-bot/solab
ExecStart=/usr/bin/python3 /home/root/telegram-bot/solab/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
